# infra/systemd/tradingbot-webhook-worker.service
[Unit]
Description=Trading Bot Telegram Webhook Worker
After=network.target tradingbot-webhook-server.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/trading-bot
EnvironmentFile=/opt/trading-bot/.env
ExecStart=/opt/trading-bot/.venv/bin/python apps/telegram_webhook/worker.py
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/trading-bot/data

[Install]
WantedBy=multi-user.target
